# packages/sysutils/system-utils/system.d/input.service
[Unit]
Description=Volume button service
After=headphones.service
Before=rocknix.target
StartLimitIntervalSec=0

[Service]
Environment=HOME=/storage
ExecStart=/usr/bin/input_sense
Restart=always
RestartSec=1

[Install]
WantedBy=multi-user.target
