# packages/sysutils/system-utils/system.d/video.service
[Unit]
Description=Video Sense
Before=rocknix.target

[Service]
Environment=HOME=/storage
ExecStart=/usr/bin/video_sense
Restart=on-failure

[Install]
WantedBy=multi-user.target
