# packages/sysutils/autostart/profile.d/002-autostart
#!/bin/bash
# SPDX-License-Identifier: GPL-2.0
# Copyright (C) 2023 JELOS (https://github.com/JustEnoughLinuxOS)

### Apply device quirks
if [ -e "/sys/firmware/devicetree/base/model" ]
then
  export QUIRK_DEVICE="$(strings /sys/firmware/devicetree/base/model)"
else
  export QUIRK_DEVICE="$(strings /sys/class/dmi/id/sys_vendor) $(strings /sys/class/dmi/id/product_name)"
fi
export QUIRK_DEVICE="$(echo ${QUIRK_DEVICE} | sed -e "s#[/]#-#g")"
