# packages/sysutils/autostart/system.d/rocknix-autostart.service
[Unit]
Description=ROCKNIX user autostart script
Before=emustation.service weston.service
After=network-base.service graphical.target userconfig.service

[Service]
Type=oneshot
Environment=HOME=/storage
ExecStart=/bin/bash -a -c 'source /etc/profile && exec /usr/bin/autostart'
RemainAfterExit=yes

[Install]
WantedBy=rocknix.target
